# rootfs/etc/s6-overlay/s6-rc.d/brother-ql-web/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Brother-QL-Web
# Runs the webui with Brother-QL Web Python script
# ==============================================================================
bashio::log.info "Starting Brother-QL-Web..."

APP_DIR="/app"
DEFAULT_FONTS_DIR="/config/fonts/"
EXAMPLE_LABEL_SIZE="62"
EXAMPLE_ORIENTATION="standard"
DEFAULT_LOGLEVEL="info"

# Retrieve configuration values
PRINTER=$(bashio::config 'printer')
MODEL=$(bashio::config 'model')
DEFAULT_LABEL_SIZE=$(bashio::config 'default_label_size')
DEFAULT_ORIENTATION=$(bashio::config 'default_orientation')
LOGLEVEL=$(bashio::config 'log_level')
# FONTS_DIR=$(bashio::config 'fonts_directory')
bashio::log.info "PRINTER=$PRINTER"
bashio::log.info "MODEL=$MODEL"
bashio::log.info "DEFAULT_LABEL_SIZE=$DEFAULT_LABEL_SIZE"
bashio::log.info "DEFAULT_ORIENTATION=$DEFAULT_ORIENTATION"
bashio::log.info "LOGLEVEL=$LOGLEVEL"

# Define additional environment variables
export GRADIO_ROOT_PATH=$(bashio::addon.ingress_entry)
bashio::log.info "Serving app from $GRADIO_ROOT_PATH"

# DEFAULT_LABEL_SIZE
if bashio::config.has_value "default_label_size"; then
    bashio::log.info "Using chosen default label size: '$DEFAULT_LABEL_SIZE'"
else
    DEFAULT_LABEL_SIZE=$EXAMPLE_LABEL_SIZE
    bashio::log.info "Using example label size = '$EXAMPLE_LABEL_SIZE'..."
fi

# DEFAULT_ORIENTATION
if bashio::config.has_value "default_orientation"; then
    bashio::log.info "Using chosen default orientation: '$DEFAULT_ORIENTATION'"
else
    DEFAULT_ORIENTATION=$EXAMPLE_ORIENTATION
    bashio::log.info "Using example orientation = '$EXAMPLE_ORIENTATION'..."
fi

# FONTS_DIR
FONTS_DIR=$DEFAULT_FONTS_DIR
mkdir -p $FONTS_DIR
bashio::log.info "Using default local storage for fonts."
chmod -R 775 $FONTS_DIR
bashio::log.info "Using custom fonts in $FONTS_DIR..."

# LOGLEVEL
if bashio::config.has_value "log_level"; then
    bashio::log.info "Using chosen log level: '$LOGLEVEL'"
else
    LOGLEVEL=$DEFAULT_LOGLEVEL
    bashio::log.info "Using default log level = '$DEFAULT_LOGLEVEL'..."
fi

# Navigate to application directory and execute the Python script
cd $APP_DIR
exec python3 -u run.py \
    --model $MODEL \
    --default-label-size $DEFAULT_LABEL_SIZE \
    --default-orientation $DEFAULT_ORIENTATION \
    --font-folder $FONTS_DIR \
    --loglevel $LOGLEVEL \
    tcp://$PRINTER:9100